# raspberry_pi/systemd/mqtt-publisher.service
# /etc/systemd/system/mqtt-publisher.service
[Unit]
Description=MQTT Publisher for Beeper Keeper Sensors
Documentation=https://github.com/z3r0-c001/beeperKeeper
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi_user
Group=pi_user
WorkingDirectory=/opt/beeperKeeper
ExecStart=/usr/bin/python3 -u /opt/beeperKeeper/mqtt_publisher.py
Restart=always
RestartSec=10

# Resource limits
# Increased CPU quota from 15% to 50% to allow ML inference (TensorFlow is CPU-intensive)
CPUQuota=50%
MemoryMax=192M
TasksMax=25

# Logging
StandardOutput=append:/opt/beeperKeeper/logs/mqtt_publisher.log
StandardError=append:/opt/beeperKeeper/logs/mqtt_publisher.log
SyslogIdentifier=mqtt-publisher

[Install]
WantedBy=multi-user.target
